# supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:nixy]
command=/bin/ash -c "/bin/nixy -f /etc/nixy/nixy.toml"
redirect_stderr=true
